# .devcontainer/s6/s6-rc.d/mcp-server/run
#!/command/with-contenv sh
# GhostDesk MCP server (HTTP on :3000).
#
# s6-rc dependencies guarantee that `sway` has been *started* before this
# service, but not that its IPC socket is ready to accept connections.
# The socket name embeds sway's PID (`sway-ipc.$UID.$PID.sock`), so we
# must discover it via glob and export SWAYSOCK before exec'ing the
# server — otherwise the first `swaymsg -t get_tree` call from the
# screen/windows tool would race and fail.
set -eu

BIN="${APP_DIR}/.venv/bin/ghostdesk"
if [ ! -x "${BIN}" ]; then
    echo "mcp-server: ${BIN} not found or not executable" >&2
    echo "mcp-server: cont-init.d/03-uv-sync.sh should have created the venv" >&2
    exit 1
fi

RUNTIME_DIR="${XDG_RUNTIME_DIR:-/run/user/1000}"
TIMEOUT=30
elapsed=0

echo "mcp-server: waiting for sway IPC socket in ${RUNTIME_DIR}..."
while [ "${elapsed}" -lt "${TIMEOUT}" ]; do
    for sock in "${RUNTIME_DIR}"/sway-ipc.*.sock; do
        if [ -S "${sock}" ]; then
            export SWAYSOCK="${sock}"
            echo "mcp-server: sway ready, SWAYSOCK=${sock}"
            exec s6-setuidgid "${RUN_USER}" "${BIN}"
        fi
    done
    sleep 1
    elapsed=$((elapsed + 1))
done

echo "mcp-server: timeout after ${TIMEOUT}s waiting for sway IPC socket" >&2
exit 1
