# docker/supervisord.conf
; ghostdesk supervisord config — shared between prod (Dockerfile.base) and
; dev (.devcontainer/Dockerfile). Single source of truth for the service
; graph. The only thing that varies between environments is the target
; user, injected via GHOSTDESK_USER:
;
;     prod: GHOSTDESK_USER=agent  (Dockerfile.base ENV)
;     dev:  GHOSTDESK_USER=vscode (.devcontainer/docker-compose.yml environment:)
;
; Everything else — commands, priorities, env expectations — is identical.
;
; pidfile and socket live under /tmp so supervisord can run unprivileged
; in dev (launched as vscode by the "GhostDesk: Start stack" task) without
; sudo. In prod supervisord still runs as PID 1 / root; /tmp is fine too.
;
; In dev, supervisord is NOT auto-started. The developer launches it via
; the "GhostDesk: Start stack" task in .vscode/tasks.json. In prod it is
; exec'd by docker/entrypoint.sh when no CMD is set (see its `$# -eq 0`
; branch).

[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ----------------------------------------------------------------------------
; dbus-session — session bus for sway, mako, portals.
; ----------------------------------------------------------------------------
[program:dbus-session]
command=dbus-daemon --session --nofork --nopidfile --address=%(ENV_DBUS_SESSION_BUS_ADDRESS)s
user=%(ENV_GHOSTDESK_USER)s
priority=10
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ----------------------------------------------------------------------------
; sway — wlroots compositor, headless backend via WLR_BACKENDS=headless.
; ----------------------------------------------------------------------------
[program:sway]
command=sway --unsupported-gpu
user=%(ENV_GHOSTDESK_USER)s
priority=20
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ----------------------------------------------------------------------------
; mako — notification daemon. Needs sway + dbus-session.
; ----------------------------------------------------------------------------
[program:mako]
command=mako
user=%(ENV_GHOSTDESK_USER)s
priority=30
autorestart=true
startretries=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ----------------------------------------------------------------------------
; wayvnc — VNC server bound to 127.0.0.1:5900, auth from
; ${HOME}/.config/wayvnc/config (rendered at boot from GHOSTDESK_VNC_PASSWORD
; by entrypoint.sh).
; ----------------------------------------------------------------------------
[program:wayvnc]
command=wayvnc
user=%(ENV_GHOSTDESK_USER)s
priority=40
autorestart=true
startretries=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ----------------------------------------------------------------------------
; websockify — noVNC bridge, :6080 → 127.0.0.1:5900. Relies on autorestart
; until wayvnc is listening. TLS is mandatory: websockify serves https://
; and wss:// using the cert at /etc/ghostdesk/tls/server.{crt,key}, which
; entrypoint.sh generates (self-signed in dev) or the operator mounts
; (real cert in prod). There is no plain-HTTP fallback, by design.
; ----------------------------------------------------------------------------
[program:websockify]
command=websockify --web=/usr/share/novnc --cert=/etc/ghostdesk/tls/server.crt --key=/etc/ghostdesk/tls/server.key 6080 127.0.0.1:5900
user=%(ENV_GHOSTDESK_USER)s
priority=50
autorestart=true
startretries=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ----------------------------------------------------------------------------
; mcp-server — GhostDesk MCP HTTP server on :3000. Wrapped in a small
; shell script that waits for sway's IPC socket and exports SWAYSOCK before
; exec'ing the real binary.
; ----------------------------------------------------------------------------
[program:mcp-server]
command=/usr/local/bin/mcp-run
user=%(ENV_GHOSTDESK_USER)s
priority=60
autorestart=true
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
