# .devcontainer/s6/s6-rc.d/xdg-desktop-portal-wlr/run
#!/command/with-contenv sh
exec s6-setuidgid vscode /usr/libexec/xdg-desktop-portal-wlr
