# .devcontainer/s6/s6-rc.d/websockify/run
#!/command/with-contenv sh
exec s6-setuidgid vscode websockify --web=/usr/share/novnc 6080 127.0.0.1:5900
