# .devcontainer/s6/s6-rc.d/dbus-session/run
#!/command/with-contenv sh
exec s6-setuidgid vscode dbus-daemon \
    --session \
    --nofork \
    --nopidfile \
    --address="${DBUS_SESSION_BUS_ADDRESS}"
