# .docker/supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:dbus]
command=dbus-daemon --session --nofork --address=unix:path=/tmp/dbus-session-bus
priority=50
autorestart=true
user=%(ENV_RUN_USER)s
environment=HOME="/home/%(ENV_RUN_USER)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension GLX +render -noreset
priority=100
autorestart=true
user=%(ENV_RUN_USER)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:openbox]
command=/usr/local/bin/start-openbox.sh
priority=200
autorestart=true
user=%(ENV_RUN_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/%(ENV_RUN_USER)s",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display %(ENV_DISPLAY)s -forever -shared -rfbport 5900 -passwd %(ENV_VNC_PASSWORD)s -noxdamage -no6
priority=300
autorestart=true
user=%(ENV_RUN_USER)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web /usr/share/novnc 6080 localhost:5900
priority=350
autorestart=true
user=%(ENV_RUN_USER)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:ghostdesk]
command=/usr/local/bin/uv run ghostdesk
directory=%(ENV_APP_DIR)s
priority=500
autorestart=true
user=%(ENV_RUN_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/%(ENV_RUN_USER)s",PATH="/usr/local/bin:/home/%(ENV_RUN_USER)s/.local/bin:%(ENV_PATH)s",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
