# .docker/supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:dbus]
command=dbus-daemon --session --nofork --address=unix:path=/tmp/dbus-session-bus
priority=50
autorestart=true
user=agent
environment=HOME="/home/agent"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:at-spi]
command=/usr/libexec/at-spi-bus-launcher
priority=150
autorestart=true
user=agent
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/agent",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension GLX +render -noreset
priority=100
autorestart=true
user=agent
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:openbox]
command=openbox
priority=200
autorestart=true
user=agent
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/agent",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display %(ENV_DISPLAY)s -forever -shared -rfbport 5900 -passwd %(ENV_VNC_PASSWORD)s -noxdamage -no6
priority=300
autorestart=true
user=agent
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web /usr/share/novnc 6080 localhost:5900
priority=350
autorestart=true
user=agent
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:ghostdesk]
command=/usr/local/bin/uv run ghostdesk
directory=/app
priority=500
autorestart=true
user=agent
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/agent",PATH="/usr/local/bin:/home/agent/.local/bin:%(ENV_PATH)s",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
