# docker/supervisor/conf.d/20-sway.conf
; wlroots compositor, headless backend via WLR_BACKENDS=headless.
[program:sway]
command=sway --unsupported-gpu
user=%(ENV_GHOSTDESK_USER)s
priority=20
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
