# docker/supervisor/conf.d/40-wayvnc.conf
; VNC server, pinned to 127.0.0.1. When TLS cert is mounted the config
; rendered by entrypoint.sh § 7 enables VeNCrypt + username/password;
; otherwise auth is disabled and wayvnc serves plain RFB on loopback
; only. The browser always reaches wayvnc through websockify on :6080.
[program:wayvnc]
command=wayvnc
user=%(ENV_GHOSTDESK_USER)s
priority=40
autorestart=true
startretries=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
