# docker/supervisor/conf.d/50-websockify.conf
; noVNC bridge, :6080 → 127.0.0.1:5900. Wrapper picks up
; /etc/ghostdesk/tls/server.{crt,key} if mounted and switches to wss://,
; otherwise serves plain ws/http.
[program:websockify]
command=/usr/local/bin/websockify-run
user=%(ENV_GHOSTDESK_USER)s
priority=50
autorestart=true
startretries=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
