# docker/supervisor/conf.d/60-mcp.conf
; GhostDesk MCP HTTP server on :3000. Wrapper waits for sway's IPC
; socket and exports SWAYSOCK before exec'ing the real binary.
[program:mcp-server]
command=/usr/local/bin/mcp-run
user=%(ENV_GHOSTDESK_USER)s
priority=60
autorestart=true
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
